# winspotlightgnome.service
[Unit]
Description=WinSpotlightGNOME
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/winspotlightgnome
